# analyzer/nat.crontab
# Keep the traffic forwarding up at 7am, bring it down at 8pm for each workday.

0 7 * * mon,tue,wed,thu,fri root /bin/bash /root/nat.sh start
0 20 * * mon,tue,wed,thu,fri root /bin/bash /root/nat.sh stop
