# analyzer/nat.crontab
# Keep the traffic forwarding up at 7am, bring it down at 7h30pm for each workday.

0 7 * * mon,tue,wed,thu,fri root /bin/bash /root/nat.sh start && /usr/bin/logger "bringing up traffic forwarding"
30 19 * * mon,tue,wed,thu,fri root /bin/bash /root/nat.sh stop && /usr/bin/logger "bringing down traffic forwarding"
