# btsync/btsync.conf
# upstart service manifest
description "BTSync Service"
start on startup
setuid btsync
setgid btsync
expect fork
exec /usr/local/sbin/btsync --config /etc/btsyncrc
